# launchd/com.coderails.routine-sweeper.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.coderails.routine-sweeper.watch</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/Users/harrison/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/harrison/Github/coderails/skills/dashboard/runner/bin/sweeper.sh</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Users/harrison/.claude/coderails-dashboard/queue</string>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/harrison/.claude/coderails-dashboard/routines/sweeper.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/harrison/.claude/coderails-dashboard/routines/sweeper.log</string>
</dict>
</plist>
